# systemd/wg-web.service
[Unit]
Description=WireGuard Web UI Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/wgweb/web
ExecStart=/usr/bin/python3 /path/to/wgweb/backend/server.py
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target